# crontab
*/10 * * * * cp $(/home/util/mbin/print_random_file.pl /home/util/Variety) /home/util/Greeter/background.jpg

