# crontab
*/10 * * * * cp $(/opt/LoginWalpaperUtil-xfce-variety/mbin/print_random_file.pl /opt/LoginWalpaperUtil-xfce-variety/Variety) /opt/LoginWalpaperUtil-xfce-variety/Greeter/background.jpg

